# deploy/systemd/paytool-report.timer
[Unit]
Description=Pay Transparency Tool — Daily report timer (06:00 UTC)

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
